# root/etc/s6-overlay/s6-rc.d/init-custom/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
#
# Optional one-shot init: install extra apt packages and/or run a user-supplied
# shell script BEFORE the runner service starts. Runs as root; the runner
# service drops to the LSIO `abc` user (uid 911) afterwards. No-op when none
# of EXTRA_PACKAGES / EXTRA_APT_REPOS / EXTRA_INIT_SCRIPT are set.
#
# Variables:
#   EXTRA_APT_REPOS   Semicolon-separated apt repository lines added to
#                     /etc/apt/sources.list.d/extra.list before installing.
#   EXTRA_PACKAGES    Space-separated apt package names; validated against
#                     ^[a-z0-9][a-z0-9+.\-]+$ to block shell injection.
#   EXTRA_INIT_SCRIPT Path to a shell script (typically bind-mounted RO)
#                     run with `set -e` after package installation.
#
# Incompatible with read-only rootfs (apt needs to write). For repeated use,
# build a custom image (FROM blackoutsecure/github-runner) instead.
#
set -uo pipefail

# shellcheck disable=SC1091
. /usr/local/bin/log-functions.sh
LOG_TAG="init-custom"

# Normalize: trim whitespace and treat whitespace-only values as unset, so a
# stripped-out compose template like "   " doesn't accidentally enable us.
_strip() { local s="${1:-}"; s="${s#"${s%%[![:space:]]*}"}"; s="${s%"${s##*[![:space:]]}"}"; printf '%s' "${s}"; }
EXTRA_PACKAGES="$(_strip "${EXTRA_PACKAGES:-}")"
EXTRA_APT_REPOS="$(_strip "${EXTRA_APT_REPOS:-}")"
EXTRA_INIT_SCRIPT="$(_strip "${EXTRA_INIT_SCRIPT:-}")"

# Fast-path: nothing to do. Exit before any apt or user-script work so the
# default deployment can't be tricked into running unexpected code here.
if [[ -z "${EXTRA_PACKAGES}" && -z "${EXTRA_APT_REPOS}" && -z "${EXTRA_INIT_SCRIPT}" ]]; then
    log "debug" "init-custom: no EXTRA_* triggers set, skipping (no-op)"
    exit 0
fi

# Repos alone, no packages, is almost certainly a misconfiguration: we'd
# write the apt sources but never install anything. Warn and skip cleanly.
if [[ -n "${EXTRA_APT_REPOS}" && -z "${EXTRA_PACKAGES}" ]]; then
    log "warn" "EXTRA_APT_REPOS is set but EXTRA_PACKAGES is empty -- apt sources will not be added"
    EXTRA_APT_REPOS=""
    if [[ -z "${EXTRA_INIT_SCRIPT}" ]]; then
        log "debug" "init-custom: nothing left to do, skipping"
        exit 0
    fi
fi

# LSIO non-root mode (https://docs.linuxserver.io/misc/non-root/): with
# `--user <uid>:<gid>` init runs unprivileged, so apt/EXTRA_INIT_SCRIPT
# root actions cannot be performed.
if [[ "${EUID:-$(id -u)}" -ne 0 ]]; then
    if [[ -n "${EXTRA_PACKAGES:-}" || -n "${EXTRA_APT_REPOS:-}" ]]; then
        log "warn" "Container is running as non-root (--user) -- EXTRA_PACKAGES/EXTRA_APT_REPOS are ignored"
        log "warn" "Build a custom image (FROM blackoutsecure/docker-github-runner + RUN apt-get install ...) instead"
    fi
    if [[ -n "${EXTRA_INIT_SCRIPT:-}" ]]; then
        log "warn" "EXTRA_INIT_SCRIPT will run as the non-root user; root-only actions inside the script will fail"
    fi
    # Fall through and still run EXTRA_INIT_SCRIPT (it can still do non-root work).
    EXTRA_PACKAGES=""
    EXTRA_APT_REPOS=""
fi

# Detect a read-only rootfs early -- apt cannot work in that case.
if [[ -n "${EXTRA_PACKAGES:-}" ]]; then
    if ! ( : > /etc/.init-custom.rwtest ) 2>/dev/null; then
        log "fatal" "EXTRA_PACKAGES is set but the container rootfs is read-only."
        log "fatal" "Either remove read_only: true from your compose file, or"
        log "fatal" "build a custom image: FROM ${IMAGE_NAME:-blackoutsecure/docker-github-runner} + RUN apt-get install ..."
        exit 1
    fi
    rm -f /etc/.init-custom.rwtest
fi

# ---------------------------------------------------------------------------
# Optional: add extra apt repositories
# ---------------------------------------------------------------------------
if [[ -n "${EXTRA_APT_REPOS:-}" ]]; then
    log "info" "Adding extra apt repositories"
    : > /etc/apt/sources.list.d/extra.list
    IFS=';' read -ra _repos <<< "${EXTRA_APT_REPOS}"
    for repo in "${_repos[@]}"; do
        repo="$(echo "${repo}" | sed 's/^[[:space:]]*//;s/[[:space:]]*$//')"
        [[ -z "${repo}" ]] && continue
        echo "${repo}" >> /etc/apt/sources.list.d/extra.list
        log "info" "  + ${repo}"
    done
fi

# ---------------------------------------------------------------------------
# Optional: install extra apt packages
# ---------------------------------------------------------------------------
if [[ -n "${EXTRA_PACKAGES:-}" ]]; then
    # Strict whitelist so an attacker-controlled EXTRA_PACKAGES cannot inject
    # shell commands.
    PKG_RE='^[a-z0-9][a-z0-9+.\-]+$'
    declare -a _safe_pkgs=()
    for pkg in ${EXTRA_PACKAGES}; do
        if [[ "${pkg}" =~ ${PKG_RE} ]]; then
            _safe_pkgs+=("${pkg}")
        else
            log "fatal" "EXTRA_PACKAGES contains an invalid package name: '${pkg}'"
            log "fatal" "Allowed pattern: ${PKG_RE} (lowercase letters, digits, +, ., -)"
            exit 1
        fi
    done

    if (( ${#_safe_pkgs[@]} == 0 )); then
        log "warn" "EXTRA_PACKAGES was set but contained no valid package names"
    else
        log "info" "Installing ${#_safe_pkgs[@]} extra package(s): ${_safe_pkgs[*]}"
        export DEBIAN_FRONTEND=noninteractive
        if ! apt-get update -o Acquire::Retries=3 >/tmp/init-custom-apt.log 2>&1; then
            log "fatal" "apt-get update failed -- see container logs"
            cat /tmp/init-custom-apt.log
            exit 1
        fi
        if ! apt-get install -y --no-install-recommends "${_safe_pkgs[@]}" \
                >>/tmp/init-custom-apt.log 2>&1; then
            log "fatal" "apt-get install failed for: ${_safe_pkgs[*]}"
            cat /tmp/init-custom-apt.log
            exit 1
        fi
        log "info" "Extra packages installed successfully"
        rm -rf /var/lib/apt/lists/* /tmp/init-custom-apt.log
    fi
fi

# ---------------------------------------------------------------------------
# Optional: run a user-provided init script
# ---------------------------------------------------------------------------
if [[ -n "${EXTRA_INIT_SCRIPT:-}" ]]; then
    if [[ ! -f "${EXTRA_INIT_SCRIPT}" ]]; then
        log "fatal" "EXTRA_INIT_SCRIPT is set to '${EXTRA_INIT_SCRIPT}' but the file does not exist"
        exit 1
    fi
    if [[ ! -r "${EXTRA_INIT_SCRIPT}" ]]; then
        log "fatal" "EXTRA_INIT_SCRIPT '${EXTRA_INIT_SCRIPT}' is not readable"
        exit 1
    fi

    # Refuse world-writable scripts: anyone with shell access to a writable
    # mount could otherwise pivot to root in the container by editing the
    # script between starts.
    _perms="$(stat -c '%a' "${EXTRA_INIT_SCRIPT}" 2>/dev/null || echo 000)"
    if [[ "${_perms: -1}" =~ [2367] ]]; then
        log "fatal" "EXTRA_INIT_SCRIPT '${EXTRA_INIT_SCRIPT}' is world-writable (mode ${_perms}); refusing to execute"
        log "fatal" "Fix on host: chmod o-w '${EXTRA_INIT_SCRIPT}'"
        exit 1
    fi

    log "info" "Running user init script: ${EXTRA_INIT_SCRIPT} (mode ${_perms})"
    # Run under bash with -e so a failing command aborts container startup.
    if ! bash -e "${EXTRA_INIT_SCRIPT}"; then
        log "fatal" "User init script failed: ${EXTRA_INIT_SCRIPT}"
        exit 1
    fi
    log "info" "User init script completed successfully"
fi

log "info" "Custom init complete"
