# root/etc/s6-overlay/s6-rc.d/init-gh-runner-config/up
/etc/s6-overlay/s6-rc.d/init-gh-runner-config/run
