# root/etc/s6-overlay/s6-rc.d/svc-gh-runner/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -uo pipefail

# Pin output to PID 1's stdout/stderr (the container's docker-log fds). s6
# tears down all longruns concurrently during shutdown; lines written to the
# inherited service pipe can be cut mid-flight and dropped from `docker logs`.
# Writing straight to /proc/1/fd/* survives teardown ordering.
if [[ -w /proc/1/fd/1 ]]; then exec 1>>/proc/1/fd/1; fi
if [[ -w /proc/1/fd/2 ]]; then exec 2>>/proc/1/fd/2; fi

# shellcheck disable=SC1091
. /usr/local/bin/log-functions.sh
LOG_TAG="svc-gh-runner"

EXIT_CODE="$1"

# Emit shutdown notice IMMEDIATELY (before any network/api work) -- some log
# consumers (e.g. balenaCloud dashboard streaming pane) close the stream as
# soon as the container starts to exit, so late lines may never reach the
# viewer even though Docker still captures them in the log file.
printf '%s svc-gh-runner[info]: SHUTDOWN signal received -- starting deregistration (exit code: %s)\n' \
    "$(date -u +'%Y-%m-%dT%H:%M:%SZ')" "${EXIT_CODE}"
sync 2>/dev/null || true

log "info" "Runner service stopped (exit code: ${EXIT_CODE}) -- running finish hook"
log "info" "Will attempt to deregister runner from GitHub before container exit"

# ---------------------------------------------------------------------------
# Deregister the runner from GitHub so it doesn't linger as "offline". Two
# strategies, in order:
#   1. API DELETE by runner id (requires GITHUB_PAT/GITHUB_TOKEN with the
#      same scope used for registration). Preferred -- single call, works
#      even when local credentials are stale.
#   2. config.sh remove with a freshly minted removal token (fallback when
#      we can't read the runner id, e.g. .runner missing).
# RUNNER_TOKEN is intentionally NOT used as fallback -- it's a registration
# token and the remove endpoint rejects it.
# ---------------------------------------------------------------------------
deregister_runner() {
    cd /opt/runner-bin

    # Skip if there's nothing configured
    if [[ ! -f ".runner" ]]; then
        log "info" "No .runner config present -- nothing to deregister"
        return 0
    fi

    # Resolve *_FILE secret variants the same way init-gh-runner-config does,
    # since exports don't cross s6 service boundaries.
    _resolve_file_secret() {
        local var="$1" file_var="${1}_FILE"
        local file_path="${!file_var:-}"
        if [[ -n "${file_path}" && -z "${!var:-}" && -r "${file_path}" ]]; then
            local v
            v="$(< "${file_path}")"
            export "${var}=${v%%$'\n'}"
        fi
    }
    _resolve_file_secret GITHUB_PAT
    _resolve_file_secret GITHUB_TOKEN
    _resolve_file_secret RUNNER_URL

    local auth_token="${GITHUB_PAT:-${GITHUB_TOKEN:-}}"
    local url_path="${RUNNER_URL#https://github.com/}"
    url_path="${url_path%/}"
    local api_base="https://api.github.com"

    # Resolve the API scope (enterprise/org/repo) once.
    local scope_path
    if [[ "${url_path}" == enterprises/* ]]; then
        scope_path="enterprises/${url_path#enterprises/}"
    elif [[ "${url_path}" == */* ]]; then
        scope_path="repos/${url_path}"
    else
        scope_path="orgs/${url_path}"
    fi

    # ----- Strategy 1: API DELETE by runner ID -----
    if [[ -n "${auth_token}" ]]; then
        local runner_id
        runner_id="$(jq -r '.agentId // empty' .runner 2>/dev/null)"
        if [[ -n "${runner_id}" ]]; then
            log "info" "Deregistering runner id=${runner_id} via API DELETE (${scope_path})..."
            local http_code
            # 8s cap bounds the worst case; happy path is sub-second.
            http_code="$(curl -s -o /tmp/.dereg_resp.$$ -w '%{http_code}' \
                --max-time 8 -X DELETE \
                -H "Authorization: token ${auth_token}" \
                -H "Accept: application/vnd.github+json" \
                "${api_base}/${scope_path}/actions/runners/${runner_id}" 2>/dev/null || echo 000)"
            local body
            body="$(cat /tmp/.dereg_resp.$$ 2>/dev/null || true)"
            rm -f /tmp/.dereg_resp.$$
            case "${http_code}" in
                204)
                    log "info" "Runner deregistered successfully (API DELETE 204)"
                    return 0
                    ;;
                404)
                    log "info" "Runner already gone from GitHub (API DELETE 404) -- nothing to do"
                    return 0
                    ;;
                401|403)
                    log "warn" "API DELETE rejected (HTTP ${http_code}) -- PAT/token lacks scope (need admin:org / repo / manage_runners:enterprise). Body: ${body}"
                    ;;
                422)
                    log "warn" "API DELETE refused (HTTP 422) -- runner is busy with a job. Body: ${body}"
                    ;;
                000)
                    log "warn" "API DELETE network failure (timeout/unreachable) -- will try config.sh fallback"
                    ;;
                *)
                    log "warn" "API DELETE returned HTTP ${http_code} -- will try config.sh fallback. Body: ${body}"
                    ;;
            esac
        else
            log "warn" "Could not read runner id from .runner -- skipping API DELETE"
        fi
    else
        log "info" "No GITHUB_PAT/GITHUB_TOKEN in env -- skipping API DELETE. The runner will linger in GitHub as 'offline' until reaped."
        log "info" "Tip: set CLEANUP_OFFLINE_RUNNERS=true with a PAT/token on the next start to clean it up automatically."
        log "info" "Tip: GITHUB_PAT_FILE works with docker/k8s/balena secrets -- you don't have to put the token in compose."
    fi

    # ----- Strategy 2: config.sh remove with a fresh removal token -----
    local remove_token=""
    if [[ -n "${auth_token}" ]]; then
        local resp
        resp="$(curl -fsSL --max-time 6 -X POST \
            -H "Authorization: token ${auth_token}" \
            -H "Accept: application/vnd.github+json" \
            "${api_base}/${scope_path}/actions/runners/remove-token" 2>&1)" \
            && remove_token="$(echo "${resp}" | jq -r '.token // empty')"
        if [[ -z "${remove_token}" ]]; then
            log "warn" "Failed to mint removal token via API: ${resp}"
        fi
    fi

    if [[ -z "${remove_token}" ]]; then
        log "warn" "No removal token available -- runner will remain visible in GitHub until it's cleaned up (use CLEANUP_OFFLINE_RUNNERS=true on next start, or remove it manually)"
        return 1
    fi

    log "info" "Deregistering runner via config.sh remove..."
    local remove_cmd
    if [[ "${RUN_AS_NONROOT:-0}" == "1" ]] || [[ "${EUID:-$(id -u)}" -ne 0 ]]; then
        remove_cmd=(./config.sh remove --token "${remove_token}")
    else
        remove_cmd=(s6-setuidgid abc ./config.sh remove --token "${remove_token}")
    fi
    # Hard cap config.sh remove (it makes its own API calls and can hang on
    # unreachable networks). 12s is enough for a happy round-trip; longer
    # just burns finish budget. SIGTERM, then SIGKILL 3s later.
    if timeout --kill-after=3s 12s "${remove_cmd[@]}" 2>&1 \
        | awk '{ if ($0 ~ /^[[:space:]]*$/) next;
                 cmd="date -u +%Y-%m-%dT%H:%M:%SZ"; cmd|getline ts; close(cmd);
                 printf "%s svc-gh-runner[deregister]: %s\n", ts, $0; fflush() }'; then
        log "info" "Runner deregistered successfully (config.sh)"
    else
        local rc=$?
        if [[ "${rc}" == "124" || "${rc}" == "137" ]]; then
            log "warn" "config.sh remove timed out after 12s (network/GitHub unreachable?) -- runner will appear offline in GitHub until cleanup"
        else
            log "warn" "config.sh remove failed (exit ${rc}) -- runner may remain in GitHub"
        fi
    fi
}

# Worst-case total budget inside deregister_runner: ~26s (8 + 6 + 12). Well
# under timeout-finish/stop_grace_period (both 30s); happy path is <1s.
deregister_runner

# Final flush + brief grace so the docker log driver actually emits the tail
# lines (especially the deregister output) before s6 tears the container down.
sync 2>/dev/null || true

if [[ "${RUNNER_EPHEMERAL:-false}" == "true" ]]; then
    log "info" "Ephemeral mode: shutting down container"
    sync 2>/dev/null || true
    # 200ms grace so docker's log driver flushes the final line before PID 1 exits.
    sleep 0.2
    exec s6-svscanctl -t /run/service
fi

if [[ "${EXIT_CODE}" != "0" ]]; then
    log "warn" "Runner exited abnormally (exit code: ${EXIT_CODE}), s6 will restart"
fi
